# deploy/beebot.conf
;Script for supervisord to run beebot on start
[program:beebot]
command=/home/pi/bin/wake_beebot.sh
autorestart=true
user=pi
environment=HOME="/home/pi",USER="pi"
