# deploy/beebot.conf
;Script for supervisord to run beebot on start
[program:beebot]
command=make --directory=/home/pi/beebot run-aws_iot-backend
autorestart=true
user=pi
environment=HOME="/home/pi",USER="pi"
